# packages/devtools/cli/init-templates/dxos-agent.service
[Unit]
Description=DXOS agent
After=network.target
After=default.target
After=pre-dxos-agent.service
Requires=pre-dxos-agent.service
# crash unit if agent crashes > 2 times in 2 hours
StartLimitIntervalSec=7200
StartLimitBurst=2

[Service]
Type=simple
EnvironmentFile=%h/.config/dx/dxos-agent.service.env
ExecStart=/usr/bin/env -S bash -c 'export AGENT_NODE_PATH=${AGENT_NODE_PATH:-$(which node)} AGENT_BIN=${AGENT_BIN:-$(which dx)}; $AGENT_NODE_PATH $AGENT_NODE_ARGS $AGENT_BIN agent start -f $AGENT_ARGS'
LimitCORE=infinity
Restart=always
ExecStartPre=/usr/bin/env -S bash -c 'if [[ -n "$EXTRA_START_PRE" ]]; then $EXTRA_START_PRE; fi'
ExecStartPre=/usr/bin/env -S bash -c '[[ -f %h/.local/share/dx/logs/agent.log ]] && mv %h/.local/share/dx/logs/agent.log %h/.local/share/dx/logs/agent.log.$$(date -Isec) || true'
# ensure this directory is created beforehand
StandardOutput=file:%h/.local/share/dx/logs/agent.log
StandardError=inherit

[Install]
WantedBy=default.target
