# packages/devtools/cli/init-templates/pre-dxos-agent.service
[Unit]
Description=Prepare system for DXOS agent

[Service]
Type=oneshot
ExecStart=-mkdir -p %h/.local/share/dx/logs
