# packages/core/agent/templates/dxos-agent.service.template
[Unit]
Description=DXOS agent {{PROFILE}}
After=network.target
After=default.target
# crash unit if agent crashes > 2 times in 2 hours
StartLimitIntervalSec=7200
StartLimitBurst=2

[Service]
Type=simple
Environment="LOG_FILTER=info"
Environment="AGENT_NODE_PATH={{NODE_PATH}}"
Environment="AGENT_BIN={{DX_PATH}}"
Environment="AGENT_NODE_ARGS={{NODE_ARGS}}"

ExecStart=/usr/bin/env -S bash -c '$AGENT_NODE_PATH $AGENT_NODE_ARGS $AGENT_BIN agent start --foreground --profile={{PROFILE}} {{OPTIONS}}'
LimitCORE=infinity
Restart=always
ExecStartPre=/usr/bin/env -S bash -c '[[ -f {{OUT_LOG}} ]] && mv {{OUT_LOG}} {{OUT_LOG}}.$$(date -Isec) || true'
# ensure this directory is created beforehand
StandardOutput=file:{{OUT_LOG}}
StandardError=file:{{ERROR_LOG}}

[Install]
WantedBy=default.target
